# imageroot/systemd/user/simplestaff.service
# ============================================
# SimpleStaff - Nethserver 8 Systemd Service
# ============================================

[Unit]
Description=SimpleStaff Application Pod
After=network.target

[Service]
Type=simple
EnvironmentFile=%S/state/environment
Restart=always
RestartSec=5

# Clean up any existing pod/containers
ExecStartPre=-/usr/bin/podman pod rm -f simplestaff-pod
ExecStartPre=-/usr/bin/podman rm -f simplestaff-db simplestaff-backend simplestaff-frontend

# Create shared Podman Pod (TCP_PORT is assigned by NS8 core)
ExecStartPre=/usr/bin/podman pod create --name simplestaff-pod -p ${TCP_PORT}:80

# Start 1. PostgreSQL Database Container
ExecStartPre=/usr/bin/podman run -d --pod simplestaff-pod --name simplestaff-db \
  -e POSTGRES_USER=simplestaff \
  -e POSTGRES_PASSWORD=${DATABASE_PASSWORD} \
  -e POSTGRES_DB=simplestaff \
  -v %S/state/postgres_data:/var/lib/postgresql/data:Z \
  ${POSTGRES_IMAGE}

# Start 2. Rust Backend API Container
ExecStartPre=/usr/bin/podman run -d --pod simplestaff-pod --name simplestaff-backend \
  --pull=always \
  -e DATABASE_URL=postgres://simplestaff:${DATABASE_PASSWORD}@127.0.0.1:5432/simplestaff \
  -e JWT_SECRET=${JWT_SECRET} \
  -e SERVER_HOST=0.0.0.0 \
  -e SERVER_PORT=8080 \
  -e SUPER_ADMIN_EMAIL=${SUPER_ADMIN_EMAIL} \
  -e SUPER_ADMIN_PASSWORD=${SUPER_ADMIN_PASSWORD} \
  ${SIMPLESTAFF_BACKEND_IMAGE}

# Start 3. Nginx Frontend Container (Foreground process for Systemd lifecycle tracking)
ExecStart=/usr/bin/podman run --pod simplestaff-pod --name simplestaff-frontend \
  --pull=always \
  ${SIMPLESTAFF_FRONTEND_IMAGE}

# Stop & Cleanup on service stop
ExecStop=/usr/bin/podman pod stop simplestaff-pod
ExecStopPost=/usr/bin/podman pod rm -f simplestaff-pod

[Install]
WantedBy=default.target
